# templates/sonar.service.erb
[Unit]
Description=Sonar <%= @version %>
After=network.target network-online.target
Wants=network-online.target

[Service]
ExecStart=<%= @installroot -%>/sonar/bin/<%= @arch.to_s -%>/sonar.sh start
ExecStop=<%= @installroot -%>/sonar/bin/<%= @arch.to_s -%>/sonar.sh stop
ExecReload=<%= @installroot -%>/sonar/bin/<%= @arch.to_s -%>/sonar.sh restart
PIDFile=<%= @home -%>/sonar.pid
Type=forking
User=<%= @user %>
LimitNOFILE=262144

[Install]
WantedBy=multi-user.target
